# jupyterlab/rootfs/etc/s6-overlay/s6-rc.d/jupyter/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# shellcheck disable=SC2191
# ==============================================================================
# Home Assistant Community Add-on: JupyterLab
# Runs Jupyter
# ==============================================================================
declare -a options

bashio::log.info 'Starting the Jupyter server...'

# Find the matching Jupyter log level
if bashio::config.exists 'log_level'; then
    case "$(bashio::string.lower "$(bashio::config 'log_level')")" in
        all|trace|debug)
            options+=(--log-level='DEBUG')
            ;;
        info|notice)
            options+=(--log-level='INFO')
            ;;
        warning)
            options+=(--log-level='WARN')
            ;;
        error)
            options+=(--log-level='ERROR')
            ;;
        fatal|off)
            options+=(--log-level='CRITICAL')
            ;;
    esac
else
    options+=(--log-level='INFO')
fi

# Change current working directory
cd /config/notebooks || bashio::exit.nok 'Failed changing working directory'

# Run Juypter Notebook server
exec jupyter lab "${options[@]}"
